# docker/omnibus/supervisord.conf
[unix_http_server]
file=/run/clio/supervisor.sock

[supervisorctl]
serverurl=unix:///run/clio/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/run/clio/supervisord.pid
childlogdir=/data/logs
# Pass all inherited env vars (set by entrypoint) through to child programs
strip_ansi=true

; ── PostgreSQL ─────────────────────────────────────────────────────────────────
; Priority 10 = started first; postgres must be ready before backend connects.
[program:postgres]
command=su-exec postgres postgres -D /data/pgdata
    -c "listen_addresses=127.0.0.1"
    -c "log_destination=stderr"
    -c "logging_collector=off"
priority=10
autorestart=true
startsecs=3
stopsignal=INT
stopasgroup=true
killasgroup=true
stdout_logfile=/data/logs/postgres.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/data/logs/postgres.log

; ── Redis ──────────────────────────────────────────────────────────────────────
[program:redis]
command=redis-server /run/clio/redis.conf
priority=20
autorestart=true
startsecs=2
stdout_logfile=/data/logs/redis.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/data/logs/redis.log

; ── Backend ────────────────────────────────────────────────────────────────────
; priority=30 starts after postgres and redis; the backend has its own
; waitForDatabase() and waitForRedis() retry loops, so startup order is handled.
[program:backend]
command=node /app/backend/server.js
directory=/app/backend
priority=30
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/data/logs/backend.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/backend.log

; ── Nginx ──────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
priority=40
autorestart=true
startsecs=2
stdout_logfile=/data/logs/nginx-access.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/data/logs/nginx-error.log
stdout_logfile_maxbytes=10MB
